# braiins-manager-agent/rootfs/etc/services.d/bma-daemon/run
#!/usr/bin/with-contenv bashio

cd /data || exit

if [ ! -f /usr/bin/bma-daemon ]; then
    bashio::log.error "bma-daemon binary not found"
    exit 1
fi

# Get the logging level from config or use "info" as default if missing
LOG_LEVEL=$(bashio::config 'logging.level' 'info')
exec RUST_LOG="${LOG_LEVEL}" /usr/bin/bma-daemon